# scripts/internlabs.conf
# 
# Example Upstart Service For InternLabs
# Should be placed in /etc/init/internlabs.conf
# 
# Example usage;
# start internlabs 
# 

description 'Internlabs'
author 'Steam Roller'
 
env NAME="InternLabs"
env LOG_FILE=/var/log/internlabs.log

env USER=ec2-user

env NODE_BIN=/usr/bin/node
env NODE_ENV="production"

env PORT=3001
env DIRECTORY=/var/www/internlabs.com.au
env APP=index.js
 
start on runlevel [23]
stop on shutdown

# Respawn in case of a crash, with default parameters
respawn
 
script
  cd $DIRECTORY
  su $USER
  # Make sure logfile exists and can be written by the user we drop privileges to
  touch $LOG_FILE
  chown $USER:$USER $LOG_FILE
  # recommended approach in case of su/sudo usage so that service does not fork
  exec sudo -u $USER PORT=$PORT NODE_ENV=$NODE_ENV $NODE_BIN $DIRECTORY/$APP >> $LOG_FILE 2>&1
end script
 
post-start script
  echo "app $NAME $NODE_ENV post-start event" >> $LOG_FILE
end script